# srcCorExar/99-corExar.rules
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="04e2", RUN+="/sbin/corExar $attr{busnum} $attr{devnum}"